# device/repka/repka4/mkosi.extra/etc/systemd/system/systemd-growfs-root.service.d/10-rootfs.conf
[Unit]
ConditionFirstBoot=true

[Install]
WantedBy=multi-user.target
